# ext/debian/mcollective/preinst
#!/bin/bash
set -e

if [ -x "/etc/init.d/mcollective" ]; then
    invoke-rc.d mcollective stop || true
fi

