# ext/debian/mcollective/prerm
#!/bin/bash
set -e

if [ "$1" = "failed-upgrade" ] ; then
    exit 1
fi

if [ -x "/etc/init.d/mcollective" ]; then
    invoke-rc.d mcollective stop || true
fi


