# ext/ec2demo/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

ruby /opt/rightscale/bin/fetchec2metadata.rb
ruby /opt/rightscale/bin/fetchec2userdata.rb

# If we're a mcollective node then start up automatically
grep -q server /var/spool/ec2/user-data.raw || /usr/local/bin/start-mcollective-demo.rb

touch /var/lock/subsys/local
